# socket-activation/pgadmin-proxy.socket
[Unit]
Description=pgadmin-socket

[Socket]
ListenStream=5555

[Install]
WantedBy=sockets.target
